# deploy/hf-demo/supervisord.conf
; P8 demo: four processes in one container (HF Space). supervisord has no readiness
; dependencies — start order is enforced inside the wrapper scripts (superset waits
; for ClickHouse, auto_bi waits for Superset); priorities only order the spawn.
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:clickhouse]
priority=10
; default config listens on localhost only — ClickHouse is never reachable from
; outside the container (nginx does not route to it)
command=clickhouse-server --config-file=/etc/clickhouse-server/config.xml
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:superset]
priority=20
command=/opt/demo/start-superset.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:autobi]
priority=30
command=/opt/demo/start-autobi.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:nginx]
priority=40
command=nginx -c /opt/demo/nginx.conf -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
